# packaging/com.jpy.vibecodingguard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jpy.vibecodingguard</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/jpy/Applications/Vibe Coding Guard.app/Contents/MacOS/VibeCodingGuard</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/jpy/Library/Logs/vibecoding-guard.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/jpy/Library/Logs/vibecoding-guard.err</string>
</dict>
</plist>
